# images/1.43/custom/crontab
*/5 * * * * /usr/local/bin/maintenance/runJobs >> /var/log/cron.log 2>&1
0 2,8,14,20 * * * /usr/local/bin/maintenance/rebuildData >> /var/log/cron.log 2>&1
30 22 * * * /usr/local/bin/maintenance/dumpBackup >> /var/log/cron.log 2>&1
