# images/1.31/custom/crontab
*/5 * * * * runJobs
0 2,8,14,20 * * * rebuildData
30 22 * * * dumpBackup
