# images/custom/crontab
*/5 * * * * sh /scripts/runJobs.sh
0 2,8,14,20 * * * sh /scripts/rebuildData.sh
30 22 * * * sh /scripts/dumpBackup.sh
